# sviluppo/infra/systemd/system/adoff-sdserver-flux2.service
[Unit]
Description=AdOff SD Server FLUX.2-klein (port 1237)
After=network.target

[Service]
Type=simple
User=mrxxx
WorkingDirectory=/home/mrxxx/sdcpp/models
ExecStart=/usr/bin/nice -n 10 /home/mrxxx/sdcpp/stable-diffusion.cpp/build/bin/sd-server --diffusion-model flux-2-klein-4b-Q8_0.gguf --vae flux2-vae.safetensors --llm qwen_3_4b.safetensors --diffusion-fa --listen-ip 0.0.0.0 --listen-port 1237 -t 28
Restart=always
RestartSec=15
MemoryMax=24G

[Install]
WantedBy=multi-user.target
