# sviluppo/infra/systemd/system/adoff-suggestions-analysis.timer
[Unit]
Description=Trigger settimanale AdOff Analisi Suggerimenti

[Timer]
# Sabato alle 10:07 (off-peak, minuto non tondo)
OnCalendar=Sat *-*-* 10:07:00
Persistent=true

[Install]
WantedBy=timers.target
