# sviluppo/infra/systemd/user/adoff-content-studio.service
[Unit]
Description=AdOff Content Studio (voce + immagini + montaggio)
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/mrxxx/Dropbox/1 Programmazione/Progetti/ChromePlugin/sviluppo/marketing/video-engine/tts-studio
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.cache/adoff-chatterbox-venv/bin/python -m uvicorn server:app --host 0.0.0.0 --port 8780 --log-level warning
Restart=always
RestartSec=4

[Install]
WantedBy=default.target
